# docker-scripts/supervisord.conf
[supervisord]
nodaemon=true

[program:httpd]
command=/usr/sbin/httpd -DFOREGROUND 
stdout_events_enabled=true
stderr_events_enabled=true
autostart=true
autorestart=true

